# cuubz-relay.service
[Unit]
Description=Cuubz Multiplayer Relay Server
After=network.target

[Service]
Type=simple
User=dadmin
WorkingDirectory=/var/www/html/server
ExecStart=/home/dadmin/.local/node-v22.22.0-linux-x64/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=MATCHMAKING_PORT=8765

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cuubz-relay

[Install]
WantedBy=multi-user.target
